# update-arch.service
[Unit]

Description=update

After=default.target

[Service]

ExecStart=~/autoUpdateLinux/updates-arch.sh

[Install]

WantedBy=default.target
