# update-debian.service
[Unit]

Description=update

After=default.target

[Service]

ExecStart=~/autoUpdateLinux/updates-debian.sh

[Install]

WantedBy=default.target
